# auto-stream.service
[Unit]
Description=Autostream ADC input
After=network.target network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
ExecStart=/root/auto-stream/owntone-auto-stream.sh
Restart=always
RestartSec=10

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10
SendSIGKILL=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=auto-stream

# Resource limits (optional)
Nice=-5
LimitNICE=-5

[Install]
WantedBy=multi-user.target
